# systemd/lb-api.service
[Unit]
Description=Load Balancer API Service
After=network.target

[Service]
Type=simple
User=lb-app
Group=lb-app
WorkingDirectory=/opt/lb-app/services/api
Environment="NODE_ENV=production"
Environment="DATABASE_URL=file:/opt/lb-app/data/lb-app.db"
Environment="API_PORT=4000"
Environment="API_HOST=0.0.0.0"
ExecStart=/usr/bin/node /opt/lb-app/services/api/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lb-api

[Install]
WantedBy=multi-user.target

