# systemd/lb-frontend.service
[Unit]
Description=Load Balancer Frontend Service
After=network.target lb-api.service
Wants=lb-api.service

[Service]
Type=simple
User=lb-app
Group=lb-app
WorkingDirectory=/opt/lb-app/frontend
Environment="NODE_ENV=production"
Environment="NEXT_PUBLIC_API_URL=http://0.0.0.0:4000/api/v1"
Environment="NEXT_PUBLIC_WS_URL=ws://0.0.0.0:4000"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lb-frontend

[Install]
WantedBy=multi-user.target

